# office_workstation/playbook/roles/common/files/rc.local
#!/bin/bash

# Enable Suspend/Hibernate Wakeup From All USB Devices
for f in /sys/bus/usb/devices/usb*/power/wakeup ; do echo enabled > $f ; done

exit 0
